# supervisord.conf
[program:stash]
command = /usr/local/Stash/bin/start-stash.sh
autorestart = true